# systemd/commuty-consumer.service
[Unit]
Description=Commuty-ED Video Consumer Service
Documentation=https://github.com/commuty-ed
After=network.target rabbitmq-server.service mysql.service
Wants=rabbitmq-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/commuty-ed
ExecStart=/bin/bash /var/www/commuty-ed/consumer_wrapper.sh
Restart=always
RestartSec=3
StandardOutput=append:/var/log/commuty/consumer.log
StandardError=append:/var/log/commuty/consumer-error.log

# ============================================
# RECURSOS - Optimizado para VPS 8GB RAM
# FFmpeg (libx264) puede usar ~1.5-2.5GB por video; 1 worker = 1 video a la vez
# ============================================

# Memoria: hasta 2.5GB para un solo job de video (comprimir + cortes + thumbnail)
MemoryMax=2560M
MemoryHigh=2G

# CPU: hasta 85% (FFmpeg es intensivo; el resto para MySQL, RabbitMQ, web)
CPUQuota=85%

# Prioridad de I/O
IOSchedulingClass=best-effort
IOSchedulingPriority=4

LimitNOFILE=65535
LimitNPROC=4096

# Videos grandes pueden tardar varios minutos
TimeoutStopSec=300

# PHP: 512M; FFmpeg se controla con FFMPEG_THREADS en config.php (default 2)
Environment=PHP_ENV=production
Environment=PHP_MEMORY_LIMIT=512M

# Reinicio inteligente
StartLimitIntervalSec=300
StartLimitBurst=5

# Protección del sistema
ProtectSystem=false
PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
