# systemd/commuty-resultado.service
[Unit]
Description=Commuty-ED Result Consumer Service
Documentation=https://github.com/commuty-ed
After=network.target rabbitmq-server.service mysql.service
Wants=rabbitmq-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/commuty-ed
ExecStart=/bin/bash /var/www/commuty-ed/resultado_wrapper.sh
Restart=always
RestartSec=3
StandardOutput=append:/var/log/commuty/resultado.log
StandardError=append:/var/log/commuty/resultado-error.log

# ============================================
# RECURSOS - VPS 8GB: este servicio solo escribe en BD (muy ligero)
# ============================================

MemoryMax=512M
MemoryHigh=384M

CPUQuota=40%

# Prioridad de I/O
IOSchedulingClass=best-effort
IOSchedulingPriority=6

# Límite de archivos abiertos
LimitNOFILE=32768

# Límite de procesos
LimitNPROC=2048

# Timeout
TimeoutStopSec=60

# Variables de entorno para PHP
Environment=PHP_ENV=production
Environment=PHP_MEMORY_LIMIT=256M

# Reinicio inteligente
StartLimitIntervalSec=300
StartLimitBurst=5

# Protección del sistema
ProtectSystem=false
PrivateTmp=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
