# systemd/commuty-resultado.service
[Unit]
Description=Commuty-ED Result Consumer Service
Documentation=https://github.com/tu-usuario/commuty-ed
After=network.target rabbitmq-server.service mysql.service
Wants=rabbitmq-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/commuty-ed
ExecStart=/usr/bin/php /var/www/commuty-ed/consumer_resultado.php
Restart=always
RestartSec=5
StandardOutput=append:/var/log/commuty/resultado.log
StandardError=append:/var/log/commuty/resultado-error.log

# Límites de recursos
MemoryLimit=128M
CPUQuota=25%

# Reiniciar si consume demasiada memoria
MemoryMax=256M

# Variables de entorno (opcional)
Environment=PHP_ENV=production

[Install]
WantedBy=multi-user.target
